# python-accesser/accesser.service
[Unit]
Description=A tool for solving SNI RST
Documentation=https://github.com/URenko/Accesser

[Service]
Type=simple
WorkingDirectory=/etc/accesser/
ExecStart=accesser
ExecStop=kill $MAINPID

[Install]
Wants=network-online.target
After=network-online.target nss-lookup.target
